# stage1/usr_from_kvm/files/sshd.socket
[Unit]
Description=OpenSSH Server Socket
DefaultDependencies=no

[Socket]
ListenStream=22
Accept=yes

[Install]
WantedBy=default.target
